# openwrt/packages/hue-factory-defaults/files/etc/uci-defaults/disabled.template
#!/bin/sh
exit 1  # It is important to fail so that the script is not removed via the overlay.
