# openwrt/packages/hue-factory-defaults/files/etc/uci-defaults/hbpl-724_remove_bridge_mode
#!/bin/sh

[ "$(uci get network.lan.type)" != "bridge" ] && exit 0

uci -q batch <<-EOF >/dev/null
        delete network.lan.type
        commit network
EOF

exit 0
