# openwrt/packages/hue-factory-defaults/files/etc/uci-defaults/remove_link_local_and_default_routes
#!/bin/sh

[ "$(uci get network.default_route)" != "route" ] && [ "$(uci get network.link_local_route)" != "route" ] && exit 0

uci -q batch <<-EOF >/dev/null
        delete network.default_route
        delete network.link_local_route
        commit network
EOF

exit 0
