# openwrt/packages/hue-factory-defaults/files/etc/uci-defaults/Z02-setup-networking-defaults
#!/bin/sh
# Copyright (C) 2015 Philips Lighting

abort() {
	echo "error: $*" >&2
	exit 1
}

resolveLanAddressWithDhcp() {
	uci batch <<EOF
set network.lan.proto='dhcp'
commit network
EOF
	return $?
}


addLinkLocalRoutes() {
	uci batch <<EOF
set network.lan.dns='208.67.222.222 208.67.220.220'
set network.default_route=route
set network.default_route.interface='lan'
set network.default_route.target='0.0.0.0'
set network.default_route.netmask='0.0.0.0'
set network.default_route.metric='99'
set network.link_local_route=route
set network.link_local_route.interface='lan'
set network.link_local_route.target='169.254.0.0'
set network.link_local_route.netmask='255.255.0.0'
set network.link_local_route.metric='99'
commit network
EOF
	return $?
}

disableDhcpServer() {
	uci batch <<EOF
set dhcp.lan.ignore='1'
commit dhcp
EOF
	return $?
}

useDnsMasqForResolvingDnsQueries() {
	ln -nsf /tmp/resolv.conf.auto /etc/resolv.conf
	return $?
}

# Ensure that the network and dhcp config files exists
[ -f /etc/config/network ] || exit 1
[ -f /etc/config/dhcp ] || exit 1


resolveLanAddressWithDhcp || abort "resolveLanAddressWithDhcp"
addLinkLocalRoutes || abort "addLinkLocalRoutes"
disableDhcpServer || abort "disableDhcpServer"
useDnsMasqForResolvingDnsQueries || abort "useDnsMasqForResolvingDnsQueries"

exit 0
