# openwrt/packages/secure-console/files/etc/init.d/secure-console
#!/bin/sh /etc/rc.common
# Copyright (C) 2015 Philips Lighting BV

START=11
STOP=97

syncShadowFileWithUbootString() {
	/lib/functions/secure-console.sh
}

boot() {
	syncShadowFileWithUbootString
}
