# openwrt/packages/hue-swupdate/files/etc/init.d/swupdate
#!/bin/sh /etc/rc.common
# Copyright (C) 2015 Philips Lighting BV

START=80
STOP=80
USE_PROCD=1

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1
SERVICE_PID_FILE=/var/run/swupdate.pid

SWUPDATE=/usr/sbin/swupdate
VAR_DIR=/var/swupdate
CONTROL_FIFO=${VAR_DIR}/start

# procd respawn settings
# RESTART_TIMEOUT = time(seconds) before next respawn
# FAIL_THRESHOLD = time(seconds) in which MAX_FAIL can be reached
# MAX_FAIL = max number of respawns within FAIL_THRESHOLD.
RESTART_TIMEOUT=10
FAIL_THRESHOLD=3600
MAX_FAIL=25

start_service() {
        mkdir -p ${VAR_DIR}
        procd_open_instance
        procd_set_param command ${SWUPDATE} -l ${CONTROL_FIFO}
        procd_set_param respawn ${FAIL_THRESHOLD} ${RESTART_TIMEOUT} ${MAX_FAIL}
        procd_set_param fail reboot
        procd_close_instance
}

stop() {
        service_stop ${SWUPDATE}
}
