# Software/openwrt/packages/secure-console/files/etc/uci-defaults/32-factory-keys
#!/bin/sh
# Copyright (C) 2015 Philips Lighting

ssh-factory-key -a
exit $?
