# Software/openwrt/packages/hue-factory-defaults/files/etc/uci-defaults/Z03-configure-homekit
#!/bin/sh
# Copyright (C) 2015 Philips Lighting

IFACE=br-lan
HOMEKIT_PIN=`fw_printenv -n homekit`
HOMEKIT_STATE_DIR=/var/homekit
HOMEKIT_PERSISTENT_DIR=/home/homekit/var
HAP_PORT=8080

abort() {
	echo "error: $*" >&2
	exit 1
}

configureMdns() {
	uci batch <<EOF
set hk_mdns.hk_mdns.netif=${IFACE}
set hk_mdns.hk_mdns.srv_dir=${HOMEKIT_STATE_DIR}/mdns-config
set hk_mdns.hk_mdns.ipc_server_file=${HOMEKIT_STATE_DIR}/mdns_ipc_server
set hk_mdns.hk_mdns.watchdog_file=${HOMEKIT_STATE_DIR}/hk_mdns_watchdog
commit hk_mdns
EOF
	return $?
}

configureHap() {
	uci batch <<EOF
set hk_hap.hk_hap.mdns_client_file=${HOMEKIT_STATE_DIR}/mdns_ipc_client
set hk_hap.hk_hap.mdns_server_file=${HOMEKIT_STATE_DIR}/mdns_ipc_server
set hk_hap.hk_hap.mdns_srv_dir=${HOMEKIT_STATE_DIR}/mdns-config
set hk_hap.hk_hap.mdns_srv_file=_hap._tcp
set hk_hap.hk_hap.keystorage_file=${HOMEKIT_PERSISTENT_DIR}/hk_keystorage
set hk_hap.hk_hap.accessory_state_file=${HOMEKIT_PERSISTENT_DIR}/hk_accessory_state_file
set hk_hap.hk_hap.device_password=${HOMEKIT_PIN}
set hk_hap.hk_hap.http_port=${HAP_PORT}
set hk_hap.hk_hap.mfi_pair_setup=1
set hk_hap.hk_hap.watchdog_file=${HOMEKIT_STATE_DIR}/hk_hap_watchdog
commit hk_hap
EOF
	return $?
}

configureMdns || abort "configureMdns"
configureHap || abort "configureHap"

exit 0
